# ansible/roles/evon-hub/files/mapper/deployment_mapper.cron
MAILTO=""
*/5 * * * * root bash -c 'source /opt/evon-hub/.env/bin/activate && /opt/evon-hub/mapper/mapper.py'
