# deployment_mapper/deployment_mapper.cron
MAILTO=""
*/5 * * * * root /opt/deployment_mapper/deployment_mapper.py

