# mapper/deployment_mapper.cron
MAILTO=""
*/5 * * * * root sleep {{ cron_offset }}; /opt/deployment_mapper/deployment_mapper.py

