# ansible/roles/evon-hub/files/cron.d/evon
# public ipv4 address change watchdog
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin
MAILTO=
* * * * * root evon --sync-pubip
